# templates/supervisord/php-fpm.conf
[program:php-fpm]
command=__PHP_FPM_COMMAND__
autorestart=true
